# scripts/systemd/forge-http.service
[Unit]
Description=Forge Starter — HTTP Server
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=forge
Group=forge
WorkingDirectory=/opt/forge-starter
Environment=PROCESS=http
EnvironmentFile=/opt/forge-starter/.env
ExecStart=/opt/forge-starter/bin/forge-starter
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forge-http

[Install]
WantedBy=multi-user.target
